# worker/systemd/novique-content-worker.service
[Unit]
Description=Novique.ai content discovery and article worker
Wants=network-online.target
After=network-online.target

[Service]
Type=oneshot
Environment=NODE_ENV=production
# User-manager PATH has no mise shims; pin Node 22 from mise for §12 self-sufficiency.
Environment=PATH=%h/.local/share/mise/installs/node/22/bin:/usr/local/bin:/usr/bin:/bin
EnvironmentFile=%h/.config/novique-worker/env
ExecStart=%h/.local/share/mise/installs/node/22/bin/npm --prefix ${NOVIQUE_WORKER_DIR} run worker:run
TimeoutStartSec=45min
NoNewPrivileges=true
PrivateTmp=true

[Install]
WantedBy=default.target
